# mqtt-dashboard/mqtt-dashboard.service
[Unit]
Description=MQTT Dashboard Web Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/mqtt-dashboard
Environment="PATH=/usr/bin:/usr/local/bin"
ExecStart=/usr/bin/python3 /path/to/mqtt-dashboard/app.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
